# templates/trigger/sprint@.service
[Unit]
Description=Sprint scan trigger for %i

[Service]
Type=oneshot
ExecStart=/bin/bash %h/repos/%i/scripts/sprint-scan.sh %h/repos/%i
