# templates/trigger/sprint@.timer
[Unit]
Description=Sprint poll timer for %i

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
